# .config/systemd/user/publish-ip.timer
[Unit]
Description=Publish external IP timer

[Timer]
OnBootSec=2m
OnUnitInactiveSec=1h

[Install]
WantedBy=multi-user.target
